# contrib/gentoo/sys-apps/knxd/files/knxd.init
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

# knxd-@SLOT@

pidfile="/var/tmp/knxd-@SLOT@.pid"
command="/usr/bin/knxd"
command_args="${KNXD_OPTS} --pid-file=${pidfile}"
start_stop_daemon_args="-u ${KNXD_USER}"
description="Daemon for accessing the EIB/KNX bus"

depend() {
	need net
	provide knxd
}
