# systemd/knxd.service.in
[Unit]
Description=KNX Daemon
After=network.target

[Service]
EnvironmentFile=@prefix@/etc/knxd.conf
ExecStart=@exec_prefix@/bin/knxd $KNXD_OPTS $KNXD_URL
User=knxd
Group=knxd

[Install]
WantedBy=multi-user.target
Also=knxd.socket
