# systemd/knxd.service.in
[Unit]
Description=KNX Daemon
After=network.target

[Service]
EnvironmentFile=@prefix@/etc/knxd.conf
ExecStart=@exec_prefix@/bin/knxd $EIBD_OPTS $EIBD_URL
# uncomment the next two lines if you want to start knxd as a regular user
# User=knxd
# Group=knxd

[Install]
WantedBy=multi-user.target
Also=knxd.socket
