# systemd/knxd.service.in
[Unit]
Description=KNX Daemon
After=network.target knxd.socket knxnetip.socket
Requires=knxd.socket knxnetip.socket

[Service]
EnvironmentFile=@sysconfdir@/knxd.conf
ExecStart=@bindir@/knxd $KNXD_OPTS
User=knxd
Group=knxd
Type=notify

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target network-online.target
Also=knxd.socket knxnetip.socket
