# systemd/knxd.service.in
[Unit]
Description=KNX Daemon
After=network.target knxd.socket

[Service]
EnvironmentFile=@sysconfdir@/knxd.conf
ExecStart=@exec_prefix@/bin/knxd $KNXD_OPTS
User=knxd
Group=knxd
Type=notify

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
Also=knxd.socket
