# systemd/knxd.socket.in
[Unit]
Description=KNX Daemon (socket)

[Socket]
ListenStream=@RUNDIR@/knxd
ListenStream=6720

[Install]
WantedBy=sockets.target
