# contrib/legacy-eibd/openwrt/files/eibd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 OpenWrt.org
START=50
STOP=50

start() {
	/usr/bin/eibd -d -u -T -i -S -c -p /var/run/eibd.pid ip:224.0.23.12
	sleep 1
	chmod a+rw /tmp/eib
}

stop() {
	kill $(cat /var/run/eibd.pid)
}

